# .config/systemd/user/shelly-notifications.service
[Unit]
Description=Shelly Notifications tray service
PartOf=default.target

[Service]
Type=simple
ExecStart=/usr/bin/shelly-notifications
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target